# service/covid_data_update.service
[Unit]
Description=CORD-19 data update service

[Service]
Type=oneshot
StandardOutput=null
ExecStart=/usr/bin/covid_data_update /data/qaservers
